# bin/sucros.conf
# Requires rootfs verification disabled
start on started system-services
stop on runlevel [06]
script
    /usr/local/bin/sucrose-daemon
end script
